# scripts/hydra-detect.service
[Unit]
Description=Hydra Detect v2.0 — Autonomous MAVLink Detection Service
After=network.target docker.service dev-ttyACM0.device
Requires=docker.service dev-video0.device dev-ttyACM0.device

[Service]
Restart=always
RestartSec=5s
ExecStart=/usr/bin/docker run --rm --runtime nvidia \
  --device /dev/video0:/dev/video0 \
  --device /dev/ttyACM0:/dev/ttyACM0 \
  -v /home/sorcc/hydra_detect/output_data:/data \
  -p 8080:8080 \
  --name hydra-detect hydra-detect:latest
ExecStop=/usr/bin/docker stop hydra-detect
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
